# tibet_app_start.conf
[supervisord]
nodaemon=true

[program:tibet_app]
command=/bin/bash -c "cd /usr/src/hello && tibet start"
